# systemd/game-update.service
[Unit]
Description=Game Update Orchestrator

[Service]
Type=oneshot
User=root
ExecStart=/srv/scripts/game-update.sh
