# systemd/game-backup.service
[Unit]
Description=Game Backup Job

[Service]
Type=oneshot
ExecStart=/srv/scripts/backup-game.sh
